# arch/etc/systemd/system/pacman-sync.timer
[Unit]
Description=Weekly pacman sync and upgrade

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
