# arch/systemd/user/cache-cleanup.timer
[Unit]
Description=Weekly user cache cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
