# arch/systemd/user/cache-cleanup.service
[Unit]
Description=Clean stale user data

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cache-cleanup
