# arch/systemd/user/tiger-cfo.service
[Unit]
Description=Tiger Web CFO supervisor — continuous fuzzing orchestrator
Documentation=https://github.com/hireAlanAyala/tiger_web
After=network-online.target
Wants=network-online.target

# Cap restart-storms — if it dies 10 times in 5 min, give up and
# require manual `systemctl --user start tiger-cfo` to investigate.
# (StartLimit* keys belong in [Unit], not [Service].)
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=%h/cfo

# Secret loaded from a 600-perm file outside the unit (so the unit
# itself stays committable / shareable). pat.env shape: KEY=VALUE.
EnvironmentFile=%h/Documents/personal/tiger_web/pat.env

# `cfo_supervisor.sh` already loops internally on transient failures;
# Restart=always is the belt to its suspenders — covers segfaults,
# OOM kills, and supervisor crashes. RestartSec matches the script's
# own sleep-on-error to avoid restart-stacking.
ExecStart=/bin/sh %h/Documents/personal/tiger_web/src/scripts/cfo_supervisor.sh
Restart=always
RestartSec=10s

# Resource hygiene. The fuzz fleet is CPU-bound and parallel; let it
# use all cores. No swap pressure expected (fuzzers are short-lived).
Nice=10

[Install]
WantedBy=default.target
