# arch/systemd/user/wifi-prefer.service
[Unit]
Description=Switch to preferred WiFi network

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wifi-prefer
