# arch/etc/systemd/system/systemd-networkd-wait-online.service.d/any.conf
# Satisfy network-online.target when ANY interface comes up, instead of waiting
# for all of them. Without this the laptop's unplugged ethernet port makes
# wait-online time out and the boot session is marked "degraded".
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any
